# src/transcoder/debian/postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/openofficed" ]; then
        update-rc.d openofficed defaults >/dev/null
        invoke-rc.d openofficed start || exit $?
fi
