# src/transcoder/init/unoconv-listen.conf
# Transcoder Unoconv Listening Service

description     "Unoconv Listener"
author          "Austin Trask <austin@artefactual.com>"

start on (net-device-up
          and local-filesystems
	  and runlevel [2345])
stop on runlevel [016]

respawn

exec /usr/bin/unoconv --listener

